# DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop hello-pi.service
systemctl disable hello-pi.service

exit 0