# src/hello-pi.service
[Unit]
Description=Hello Pi Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hello-pi
Group=hello-pi

StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

WorkingDirectory=/opt/hello-pi
ExecStart=/usr/bin/python3 /opt/hello-pi/src/worker.py

# Creates /var/lib/hello-pi owned by hello-pi automatically:
StateDirectory=hello-pi

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target